# rc.local
#!/bin/sh -e
if [ -e "/boot/interfaces" ]; then
su - kiosk -c 'xinit' &
fi
exit 0
